# pi/kiosk/gaia-display-blank.service
[Unit]
Description=GAIA display blank - DSI spento al boot finche kiosk o screen non lo accendono
Before=gaia-kiosk.service gaia-screen.service

[Service]
Type=oneshot
ExecStart=/home/asemico/gaia/kiosk/display_power.sh off

[Install]
WantedBy=multi-user.target
