# pi/herbarium/gaia-herbarium-mediapipe.service
[Unit]
Description=GAIA Herbarium — sorgente MediaPipe (la stanza suona in risposta a chi la abita)
After=network-online.target sound.target
Conflicts=gaia-herbarium-sim.service
# Gestito da gaia-agent — NON abilitare con systemctl enable

[Service]
Type=simple
User=asemico
WorkingDirectory=/home/asemico/gaia/herbarium
EnvironmentFile=/etc/gaia/device.conf
EnvironmentFile=-/etc/gaia/herbarium.conf
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/asemico/gaia/herbarium/mediapipe_source.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
