# minipc/camera/gaia-camera.service
[Unit]
Description=GAIA Camera Server (miniPC) — shared memory + MJPEG stream :8766
After=network.target

[Service]
Type=simple
User=core
WorkingDirectory=/home/core/core-node-0/minipc/camera
ExecStart=/media/core/D/face-env/bin/python3 /home/core/core-node-0/minipc/camera/camera_server.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

Environment=PYTHONUNBUFFERED=1
Environment=HOME=/home/core
Environment=MJPEG_PORT=8766
Environment=MJPEG_QUALITY=72
# Override camera: se webcam non è /dev/video0, impostare CAMERA_INDEX in
# /etc/gaia/camera.conf (non qui — vedi config.py, env > file > default).
# Capitato 2026-07-04: riconnessione USB ha rinumerato la webcam UC W20 da
# /dev/video0 a /dev/video1 -> gaia-camera in crash-loop ("Camera 0 non
# accessibile"), yolo/mediapipe in attesa (comportamento corretto, non il
# problema). Diagnosi: `ls /dev/video*` + provare gli indici con OpenCV
# (`cv2.VideoCapture(i).isOpened()`) per trovare quello vero.

[Install]
WantedBy=multi-user.target
