# pi/agent/gaia-voice.service
[Unit]
Description=GAIA Voice Node
After=network.target sound.target
# Gestito da gaia-agent — NON abilitare con systemctl enable

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/gaia/voice-node
ExecStart=/opt/gaia/voice-node/start.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=MQTT_HOST=192.168.1.142

[Install]
